# launchd/com.openclass.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclass.web</string>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>__LAUNCH_DIR__</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd "__LAUNCH_DIR__" &amp;&amp; export NEXT_TELEMETRY_DISABLED=1 NODE_ENV=production; while true; do if [[ ! -f "__LAUNCH_DIR__/apps/web/.next/BUILD_ID" ]]; then npm run build:web; fi; pids="$(lsof -tiTCP:3000 -sTCP:LISTEN 2&gt;/dev/null || true)"; if [[ -n "$pids" ]]; then echo "[openclass-web] $(date '+%Y-%m-%d %H:%M:%S') stopping stale process on :3000: $pids"; print -r -- "$pids" | xargs kill &gt;/dev/null 2&gt;&amp;1 || true; sleep 1; pids="$(lsof -tiTCP:3000 -sTCP:LISTEN 2&gt;/dev/null || true)"; if [[ -n "$pids" ]]; then echo "[openclass-web] $(date '+%Y-%m-%d %H:%M:%S') force-stopping stale process on :3000: $pids"; print -r -- "$pids" | xargs kill -9 &gt;/dev/null 2&gt;&amp;1 || true; fi; fi; echo "[openclass-web] $(date '+%Y-%m-%d %H:%M:%S') starting Next.js production server on :3000"; npm --prefix apps/web run start; code=$?; echo "[openclass-web] $(date '+%Y-%m-%d %H:%M:%S') server exited with code $code; restarting in 5s"; sleep 5; done</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/openclass-web.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/openclass-web.log</string>
</dict>
</plist>
